# Recipes/led/files/led.service
[Unit]
Description="this is my service"

[Service]
ExecStart=/usr/bin/led
Restart=always
Type=simple

[Install]
WantedBy=multi.user.target
